# dist/systemd/mirrorcache-backstage.service
[Unit]
Description=MirrorCache daemon for background jobs
After=postgresql.service

[Service]
User=mirrorcache
Group=mirrorcache
ExecStart=/usr/share/mirrorcache/script/mirrorcache-backstage
Nice=19
Restart=on-failure
EnvironmentFile=/etc/mirrorcache/conf.env

[Install]
WantedBy=multi-user.target
