# dist/systemd/mirrorcache-backstage.service
[Unit]
Description=MirrorCache daemon for background jobs
After=systemd-tmpfiles-setup.service postgresql.service mariadb.service

[Service]
User=mirrorcache
Group=mirrorcache
ExecStart=/usr/share/mirrorcache/script/mirrorcache-backstage
Nice=19
Restart=on-failure
RestartSec=10
EnvironmentFile=/etc/mirrorcache/conf.env
WorkingDirectory=/var/lib/mirrorcache
Environment="MOJO_TMPDIR=/tmp/mirrorcache"

[Install]
WantedBy=multi-user.target
