# dist/systemd/mirrorcache-backstage-hashes.service
[Unit]
Description=MirrorCache daemon for file hashes calculation jobs

[Service]
User=mirrorcache
Group=mirrorcache
ExecStart=/usr/share/mirrorcache/script/mirrorcache-backstage-hashes
Nice=19
Restart=on-failure
EnvironmentFile=/etc/mirrorcache/conf.env

[Install]
WantedBy=multi-user.target
