# dist/systemd/mirrorcache.service
[Unit]
Description=MirrorCache webApp
After=systemd-tmpfiles-setup.service postgresql.service mariadb.service nss-lookup.target
Before=apache2.service

[Service]
User=mirrorcache
Group=mirrorcache
ExecStart=/usr/share/mirrorcache/script/mirrorcache-daemon
EnvironmentFile=/etc/mirrorcache/conf.env
WorkingDirectory=/var/lib/mirrorcache
Environment="MOJO_TMPDIR=/tmp/mirrorcache"

[Install]
WantedBy=multi-user.target
