# dist/systemd/mirrorcache.service
[Unit]
Description=MirrorCache webApp
Before=apache2.service
After=postgresql.service nss-lookup.target
Requires=postgresql.service

[Service]
User=mirrorcache
ExecStart=/usr/share/mirrorcache/script/mirrorcache-daemon
EnvironmentFile=/usr/share/mirrorcache/conf.env

[Install]
WantedBy=multi-user.target
