# dist/systemd/mirrorcache.service
[Unit]
Description=MirrorCache webApp
Before=apache2.service
After=postgresql.service mariadb.service nss-lookup.target

[Service]
User=mirrorcache
Group=mirrorcache
ExecStart=/usr/share/mirrorcache/script/mirrorcache-daemon
EnvironmentFile=/etc/mirrorcache/conf.env
WorkingDirectory=/var/lib/mirrorcache

[Install]
WantedBy=multi-user.target
