# dist/systemd/mirrorcache-subtree.service
[Unit]
Description=MirrorCache subtree webApp

[Service]
User=mirrorcache
Group=mirrorcache
ExecStart=/usr/share/mirrorcache/script/mirrorcache-daemon
EnvironmentFile=/etc/mirrorcache/conf-subtree.env

[Install]
WantedBy=multi-user.target
